# overlay/System/Library/LaunchDaemons/com.apple.notifyd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<false/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.notifyd</string>
	<!-- RunAtLoad added for FreeBSD port (Phase J runtime
	     validation). Force eager start until launchd-842 MachServices
	     lazy-spawn is verified working. -->
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>ASL_DISABLE</key>
		<string>1</string>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.system.notification_center</key>
		<true/>
	</dict>
	<key>PublishesEvents</key>
	<string>com.apple.notifyd.matching</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/sbin/notifyd</string>
	</array>
	<!-- Capture stderr for post-mortem (Phase J runtime). -->
	<key>StandardErrorPath</key>
	<string>/var/log/notifyd.stderr</string>
	<key>StandardOutPath</key>
	<string>/var/log/notifyd.stderr</string>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>LimitLoadToSessionType</key>
	<string>System</string>
</dict>
</plist>
