# overlay/System/Library/LaunchDaemons/com.apple.kextd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<!-- kextd, the in-kernel IOKit matcher's loader daemon (K3b, #217).
	     Early, persistent System daemon — Apple's shape (kextd was the
	     HOST_KEXTD_PORT server). It registers HOST_KEXTD_PORT so the kernel
	     matcher can reach it, pushes the IOKitPersonalities (push-triggers-
	     match #216 then fires load requests for already-unmatched devices),
	     and serves load requests forever (OSKextLoad -> kldload -> attach).

	     Ordering is irrelevant: launchd has none, and push-triggers-match
	     makes the kernel re-fire requests when kextd pushes, so kextd may
	     start at any time without losing loads. -->
	<key>Label</key>
	<string>com.apple.kextd</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/kextd</string>
		<string>-w</string>
		<!-- -v: log the boot daemon's startup push + each load request it
		     serves, so CI can see whether boot-time autoload actually fires.
		     It must work with no manual poke (real-hardware shape). -->
		<string>-v</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<!-- Quiet by default (no -v): logs only meaningful events (loaded
	     <bundle>, errors) to a file, like other System daemons. The full
	     startup trace is available by adding -v. (Routing to syslog is a
	     later refinement — #211.) -->
	<key>StandardErrorPath</key>
	<string>/var/log/kextd.log</string>
	<key>StandardOutPath</key>
	<string>/var/log/kextd.log</string>
	<key>LimitLoadToSessionType</key>
	<string>System</string>
</dict>
</plist>
