# overlay/System/Library/LaunchDaemons/com.apple.hostnamed.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.hostnamed</string>

	<!-- One-shot at boot. iter 1 synthesizes the hostname, publishes
	     it to SCDynamicStore + kernel + notifyd, then exits. No
	     persistent loop, no MachServices (no RPC surface yet — iter 3+
	     adds SCPreferences integration). KeepAlive=false so launchd
	     doesn't keep relaunching the one-shot.

	     PAM port iter 4 (issue #99): RunAtLoad RESTORED. The
	     original iter-1 drop was specifically because the additional
	     RunAtLoad plist accelerated launchd's boot-time plist
	     dispatch enough to let getty/login race ahead of pam_xdg's
	     /var/run/xdg setup, breaking pam_open_session and wedging
	     root login. PAM port iter 3 (PR #98) dropped pam_xdg from
	     the overlay /etc/pam.d/* entirely; nothing references it
	     anymore. The race is structurally resolved.

	     With RunAtLoad=true, hostnamed runs at boot before getty
	     prints its banner; the synthesized hostname appears in
	     "FreeBSD/amd64 (HOSTNAME) (console)" instead of the
	     unset-placeholder "Amnesiac". run.sh's manual hostnamed
	     invocations remain (covering the SCPrefs + DHCP fixture
	     rounds in run.sh — those run after the boot-time one).

	     NOTE: com.apple.syslogd.plist has a separate, unrelated
	     reason its RunAtLoad is dropped (the launch_msg LAUNCH_KEY_CHECKIN
	     Mach IPC hang from the mach.ko interruptibility gap). That
	     issue isn't fixed by the PAM port — syslogd stays dropped.

	     Plan: https://pkgdemon.github.io/freebsd-hostnamed-plan.html
	     Issue: #63 (hostnamed iter 1), #99 (RunAtLoad restoration) -->
	<key>ProgramArguments</key>
	<array>
		<string>/usr/sbin/hostnamed</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>

	<!-- Capture stderr for post-mortem + the HOSTNAMED-OK marker
	     visibility from the on-image test runner. -->
	<key>StandardErrorPath</key>
	<string>/var/log/hostnamed.stderr</string>
	<key>StandardOutPath</key>
	<string>/var/log/hostnamed.stderr</string>

	<!-- System session — runs as root before any user logs in. Same
	     scope as configd / hwregd / notifyd / mDNSResponder. -->
	<key>LimitLoadToSessionType</key>
	<string>System</string>
</dict>
</plist>
