# overlay/System/Library/LaunchDaemons/com.apple.mDNSResponder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.apple.mDNSResponder.plist

  Boot-time launchd plist for mDNSResponder (freebsd-launchd-mach
  Bonjour/zeroconf daemon).

  iter 1 is the daemon skeleton: mDNSResponder claims the
  com.apple.mDNSResponder Mach service via bootstrap_check_in and
  sleeps until SIGTERM. iter 2 vendors mDNSCore + mDNSPosix +
  mDNSShared and wires the multicast engine + the dns_sd.h client
  wire protocol over Mach. RunAtLoad + KeepAlive so it starts at
  PID-1 launchd boot and respawns on crash. StandardErrorPath /
  StandardOutPath route the daemon's log to
  /var/log/mDNSResponder.stderr — same convention hwregd / configd /
  ipconfigd use; the CI boot test dumps the file post-boot for
  verification.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.apple.mDNSResponder</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/sbin/mDNSResponder</string>
        <!--
          -debug: tell PosixDaemon.c's ParseCmdLineArgs to set
          mDNS_DebugMode. This does two things:

          (1) skips daemon(0, 0) just like -foreground did — without
              that, daemon() closes stdin/stdout/stderr and we lose all
              fprintf-to-stderr output (including the iter-2
              MDNS-ENGINE-OK marker). launchd already manages the
              daemon's lifecycle via KeepAlive, so self-backgrounding is
              redundant anyway.

          (2) routes mDNSPlatformWriteLogMsg through fprintf(stderr)
              instead of syslog(). The engine logs at NOTICE volume
              (mDNS_LoggingEnabled is on), and feeding that into the
              shared asl/syslog pipeline leaks the chatter onto
              /dev/console at boot, papering over the getty login
              prompt. With -debug it goes straight to StandardErrorPath
              (/var/log/mDNSResponder.stderr) — mDNS keeps its full log
              in its own file, off the console and off syslogd.

          Note: PosixDaemon.c only inspects argv[1], so -debug must be
          the sole argument (it supersedes -foreground).
        -->
        <string>-debug</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>MachServices</key>
    <dict>
        <key>com.apple.mDNSResponder</key>
        <true/>
    </dict>
    <key>StandardOutPath</key>
    <string>/var/log/mDNSResponder.stderr</string>
    <key>StandardErrorPath</key>
    <string>/var/log/mDNSResponder.stderr</string>
</dict>
</plist>
